# etc/supervisord.d/graphite-web.ini
[program:graphite-webapp]
user = www-data
directory = /usr/share/graphite-web/webapp
environment = PYTHONPATH='/usr/share/graphite-web/webapp'
command = /usr/bin/gunicorn_django -b127.0.0.1:8000 -w2 graphite/settings.py
redirect_stderr=true
stdout_syslog = true
autorestart = true
